# projects/MSVSMd/msvs
#!/sbin/openrc-run

description="Starts the /www HTTP server."

depend() {
	need net
	use dns logger
	want nginx
}

start() {
        /usr/bin/python /scripts/MSVS.py
	nohup /usr/bin/python -m http.server --directory /www 80 &
}
